# init.d/radarr
#!/sbin/openrc-run

name="radarr"
pidfile="/run/${RC_SVCNAME}.pid"
directory="/usr/lib/radarr/bin/"
command="/usr/lib/radarr/bin/Radarr"
command_args="-nobrowser -data=/var/lib/radarr"
command_background=true

depend() {
	need net
} 
