# init.d/scanservjs
#!/sbin/openrc-run

name="scanservjs"
pidfile="/run/${RC_SVCNAME}.pid"
command="/usr/bin/node"
directory="/var/www/scanservjs/"
command_args="./server/server.js"
command_background=true

depend() {
	need net
}
