# init.d/scanservjs
#!/sbin/openrc-run

name="scanservjs"
command="/usr/bin/node"
command_args=/var/www/scanservjs/server/server.js"
command_background=true

depend() {
	need net
}
